# saved-unused/10-image-cayo/mkosi.extra/usr/lib/systemd/system/docker.socket.d/override.conf
[Unit]
ConditionKernelCommandLine=!snow-linux.live=1
