# saved-unused/10-image-cayo/mkosi.extra/usr/lib/systemd/system/user@.service.d/shutdown-order.conf
[Unit]
# Stop user services before trying to unmount user tmpfs
Before=umount.target