# shared/snow/tree/usr/lib/systemd/system/incus-startup.service.d/shutdown-order.conf
[Unit]
# Ensure incus-startup stops (which stops VMs/containers) before mounts are unmounted
Before=umount.target

[Service]
# Give more time for VMs to shut down gracefully
TimeoutStopSec=600