# proxmox-bootc-main/system_files/usr/lib/systemd/system/media.mount
[Unit]
DefaultDependencies=no
Before=local-fs.target sysinit.target

[Mount]
What=/run/media
Where=/media
Options=bind

[Install]
WantedBy=local-fs.target