# shared/cayo/tree/usr/lib/systemd/system/usr-local.mount
[Unit]
DefaultDependencies=no
Before=local-fs.target sysinit.target
After=var.mount

[Mount]
What=/var/usrlocal
Where=/usr/local
Options=bind

[Install]
WantedBy=local-fs.target