# shared/snow/tree/usr/lib/systemd/system/snap.mount
[Unit]
Before=snapd.socket

[Mount]
What=/var/lib/snapd/snap
Where=/snap
Options=bind
Type=none

[Install]
WantedBy=snapd.socket