# saved-unused/10-image-cayo/mkosi.extra/usr/lib/systemd/system/incus.service.d/shutdown-order.conf
[Unit]
# Ensure incusd stops before mount points are unmounted
Before=umount.target