# shared/packages/azurevpn/tree/usr/lib/systemd/system/microsoft-azurevpn-workaround.service
[Unit]
Description=Workaround AzureVPN not having the correct caps
ConditionPathExists=/usr/lib/microsoft-azurevpnclient/microsoft-azurevpnclient
ConditionPathExists=/usr/lib/microsoft-azurevpnclient/lib
ConditionPathExists=/usr/lib/microsoft-azurevpnclient/data
After=local-fs.target

[Service]
Type=oneshot
# Copy if it doesn't exist
ExecStartPre=/usr/bin/bash -c "[ -d /usr/local/lib/microsoft-azurevpnclient ] || /usr/bin/mkdir /usr/local/lib/microsoft-azurevpnclient"
ExecStartPre=/usr/bin/bash -c "[ -x /usr/local/lib/microsoft-azurevpnclient/.microsoft-azurevpnclient ] || /usr/bin/cp /usr/lib/microsoft-azurevpnclient/microsoft-azurevpnclient /usr/local/lib/microsoft-azurevpnclient/.microsoft-azurevpnclient"
ExecStartPre=/usr/bin/bash -c "[ -d /usr/local/lib/microsoft-azurevpnclient/lib ] || /usr/bin/ln -s /usr/lib/microsoft-azurevpnclient/lib /usr/local/lib/microsoft-azurevpnclient/lib"
ExecStartPre=/usr/bin/bash -c "[ -d /usr/local/lib/microsoft-azurevpnclient/data ] || /usr/bin/ln -s /usr/lib/microsoft-azurevpnclient/data /usr/local/lib/microsoft-azurevpnclient/data"
# This is faster than using .mount unit. Also allows for the previous line/cleanup
ExecStartPre=/usr/bin/mount --bind /usr/local/lib/microsoft-azurevpnclient/.microsoft-azurevpnclient /usr/lib/microsoft-azurevpnclient/microsoft-azurevpnclient
# Fix Caps
ExecStart=/usr/sbin/setcap cap_net_admin+eip /usr/lib/microsoft-azurevpnclient/microsoft-azurevpnclient
# Clean-up after ourselves
ExecStop=/usr/bin/umount /usr/lib/microsoft-azurevpnclient/microsoft-azurevpnclient
ExecStop=/usr/bin/rm /usr/local/lib/microsoft-azurevpnclient/.microsoft-azurevpnclient
ExecStop=/usr/bin/rm /usr/local/lib/microsoft-azurevpnclient/lib
ExecStop=/usr/bin/rm /usr/local/lib/microsoft-azurevpnclient/data
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
