# system_files/usr/lib/systemd/system/mnt.mount
[Unit]
DefaultDependencies=no
Before=local-fs.target sysinit.target

[Mount]
What=/var/mnt
Where=/mnt
Options=bind

[Install]
WantedBy=local-fs.target