# shared/cayo/tree/usr/lib/systemd/system/incus.socket.d/override.conf
[Unit]
ConditionKernelCommandLine=!cayo-linux.live=1
